# src/roles/vm/files/data.mount
[Unit]
Description=Data drive

[Mount]
What=/dev/vdb5
Where=/data
Type=xfs
Options=defaults

[Install]
WantedBy=multi-user.target